# meta-beagleboard-extras/recipes-support/usb-gadget/gadget-init/bone-gmass-eject.rules
SUBSYSTEM=="block",ACTION=="change",KERNEL=="mmcblk0p1",RUN+="/usr/bin/g-ether-start-service.sh"
